# packaging/systemd/mackesd.service.d/90-stop-policy.conf
[Service]
# REL-2 — Fedora ships a global service.d/10-timeout-abort.conf that can turn
# stop timeouts into SIGABRT+coredumps. Keep mackesd promotions non-aborting even
# when that global drop-in is present.
TimeoutStopSec=90
TimeoutStopFailureMode=terminate
