# packaging/systemd/nebula.service.d/10-mesh-recovery.conf
# MCNF drop-in for the stock nebula.service (from the `nebula` RPM):
# auto-restart the overlay on failure and notify if it gives up. The found/join
# verbs + the RPM enable nebula.service so the overlay also comes up at boot
# (the stock unit ships disabled).
[Unit]
OnFailure=mesh-alert@%n.service
StartLimitIntervalSec=300
StartLimitBurst=10

# Restart=always (not on-failure): mackesd's nebula_supervisor reloads nebula
# via SIGHUP, and this nebula build *exits* on HUP — which systemd counts as a
# CLEAN stop (SIGHUP is in the default success-signal set), so on-failure would
# NOT restart it and the overlay would stay down. `always` brings it back in
# RestartSec regardless of how it stopped. The mesh-health watchdog is the
# slower backstop.
[Service]
Restart=always
RestartSec=3
