# packaging/systemd/mackesd.service
[Unit]
Description=MCNF control-plane daemon (mackesd)
Documentation=https://github.com/matthewmackes/magic-mesh
After=network-online.target
Wants=network-online.target
# Service-manager policy: keep retrying across a generous window (a fixed role
# / mount brings it up without manual help) but don't tight-loop forever; the
# mesh-health watchdog is the backstop if the start limit is ever hit.
StartLimitIntervalSec=300
StartLimitBurst=10
# Notify when the daemon enters a failed state (the node "goes offline").
OnFailure=mesh-alert@%n.service

[Service]
# BULLETPROOF-2 — Type=notify + a watchdog so a HUNG daemon (process alive but
# the runtime wedged, e.g. a thread stuck in D-state on a bad FUSE read) is
# auto-restarted; the plain process-alive check can't see a hang. mackesd sends
# READY=1 once the supervisor + all responders are up, then pings WATCHDOG=1
# from its main loop. NotifyAccess=main: only the main process's notifications
# count. TimeoutStartSec is generous so a slow-but-not-hung init (cold mount,
# debug keygen) is never killed before READY=1; WatchdogSec restarts a true hang.
Type=notify
NotifyAccess=main
TimeoutStartSec=120
# WIFI-WATCHDOG (2026-06-29): 180 s, not 60 s. A node on a slow/high-latency link
# (a laptop on WiFi) can legitimately starve the liveness beat for tens of seconds
# during a blocking coordination op; at 60 s that false-restarted mackesd into a
# crash-loop and made in-place upgrades on WiFi impossible. The beat-freshness gate
# (MACKESD_WATCHDOG_FRESH_MS, default 60 s) + this window still catch a genuine
# wedge within ~3-4 min. Raise both together for a worse link.
WatchdogSec=180
# SETUP-fix — SHARED bus spool so this root daemon's responders + the uid-1000
# desktop GUIs (workbench/applet) share ONE bus (per-HOME data_dir() otherwise
# splits /root from /home/<u>, so every workbench<->mackesd RPC times out). The
# user session gets the same value via /usr/lib/environment.d/90-mde-bus.conf;
# the sticky 1777 dir is created by tmpfiles (magic-mesh.conf).
Environment=MDE_BUS_ROOT=/run/mde-bus
# SETUP-fix — SHARED mesh-storage dir so this root daemon and the uid-1000
# desktop GUIs resolve the SAME path. The code default is ~/QNM-Shared
# (HOME-relative), which splits /root/QNM-Shared (root daemon) from
# /home/<u>/QNM-Shared (workbench) — and /root is 0750 so the GUI can't even
# traverse it. Pinning a neutral 0755 dir (SUBSTRATE-V2: a plain Syncthing-
# replicated directory, no FUSE) lets both uids read the one shared dir. The
# user session gets the same value via environment.d/90-mde-workgroup.conf.
Environment=MDE_WORKGROUP_ROOT=/mnt/mesh-storage
ExecStart=/usr/bin/mackesd serve
# ENT-2: serve fails closed without a pinned role — that's a
# configuration error, not a crash loop; on-failure still retries
# so a fixed role brings the daemon up without manual intervention.
Restart=on-failure
RestartSec=3
# ENT-6 acceptance: kill -9 -> restarted within seconds.
TimeoutStopSec=20

# EFF-15 — the CA-backup passphrase is an OPTIONAL systemd-encrypted credential
# (ENT-11/C8): TPM/host-bound, invisible to `systemctl show` + /proc/<pid>/environ.
# It is NOT in the base unit (a `LoadCredentialEncrypted` whose source file is
# absent makes systemd refuse to start the unit — a fresh install would be
# dead-on-arrival). Instead, on a CA-holding box, provision it once and add the
# drop-in that loads it (mackesd reads $CREDENTIALS_DIRECTORY/backup-passphrase;
# SEC-7 alerts loudly if it is unset on a CA holder):
#   systemd-creds encrypt --name=backup-passphrase - \
#     /etc/mackesd/backup-passphrase.cred  < phrase.txt
#   mkdir -p /etc/systemd/system/mackesd.service.d
#   printf '[Service]\nLoadCredentialEncrypted=backup-passphrase:/etc/mackesd/backup-passphrase.cred\n' \
#     > /etc/systemd/system/mackesd.service.d/backup-passphrase.conf
#   systemctl daemon-reload && systemctl restart mackesd

# EFF-29 — service identity: root is INTENTIONAL, here is why.
# mackesd is the node's system-management plane, not a network service that
# could hide behind DynamicUser. Its role-gated workers, by design:
#   * materialize /etc/nebula/ and write /etc/kamailio-mde, /etc/rtpengine-mde
#   * merge mesh hosts into /etc/hosts (mesh_dns) + drive resolvectl
#   * drive firewalld (zones/rich rules) and systemctl restart/reload units
#   * manage CUPS queues (lpadmin), sshd overlay binds, dnf mirror pulls
# DynamicUser=yes (or a %pre system user) breaks every one of those paths or
# re-grants them via a sudo/polkit matrix wider than root itself. The trade is
# documented; containment comes from the sandbox below + dropping the cap
# classes the daemon can never legitimately need:
CapabilityBoundingSet=~CAP_SYS_MODULE CAP_SYS_RAWIO CAP_SYS_BOOT CAP_SYS_TIME
CapabilityBoundingSet=~CAP_AUDIT_CONTROL CAP_AUDIT_READ CAP_AUDIT_WRITE
CapabilityBoundingSet=~CAP_MAC_ADMIN CAP_MAC_OVERRIDE
CapabilityBoundingSet=~CAP_BLOCK_SUSPEND CAP_WAKE_ALARM CAP_LEASE CAP_LINUX_IMMUTABLE
CapabilityBoundingSet=~CAP_SYS_PACCT CAP_SYS_TTY_CONFIG CAP_IPC_LOCK
ProtectKernelModules=true
ProtectClock=true
SystemCallArchitectures=native

# EFF-12 — sandboxing. Conservative for a root network daemon that materializes
# /etc/nebula, writes /var/lib/mackesd, and shells ansible/nebula: protect the
# OS image + kernel surfaces without restricting the network/fs/subprocess paths
# it needs. (Stricter ProtectSystem/ProtectHome/RestrictAddressFamilies are
# deferred to EFF-12 follow-up — they need a live Fedora bring-up to verify
# nebula tun + the shared-storage dir keep working.)
# SELINUX-1 (2026-06-15) — NoNewPrivileges REMOVED. mackesd is the node's
# system-orchestrator: it execs helpers that legitimately change SELinux domain
# (podman → container_runtime_t for the compute-descriptor probe; libvirt/virsh;
# firewall-cmd; nebula-cert). NoNewPrivileges=yes makes the kernel forbid every
# such domain transition (`nnp_transition` denials), which broke the
# container/VM probe on a confined host AND flooded setroubleshoot with alerts
# (UNIT-EAGLE: 13×+ tokio-rt-worker → container_runtime_t, plus the libvirt
# rpc-virtqemud noise). NNP directly contradicts this unit's stated intent
# ("without restricting the network/fs/subprocess paths it needs"). Containment
# stays via the CapabilityBoundingSet drops + the Protect*/Restrict* knobs
# below, none of which block domain transitions. Do NOT re-add NoNewPrivileges.
# SUBAUDIT-D (2026-06-14) — `full` makes /etc READ-ONLY, which silently broke
# every /etc write this daemon legitimately makes (mesh_dns → /etc/hosts for
# <host>.mesh resolution; /etc/nebula materialization; /etc/ssh drop-ins).
# `true` keeps the OS image (/usr,/boot) read-only — the stated intent above —
# while leaving /etc writable for the config-management role.
ProtectSystem=true
PrivateTmp=true
# AUDIT-MESH-15 (2026-06-16) — PrivateTmp puts the daemon in a private mount
# namespace, so the FPG-7 `mount --bind` of the desktop user's XDG dirs
# (~/Documents → the communal mesh-storage dir) was invisible to the user's
# session. MountFlags=shared propagates the daemon's mounts out to the host
# namespace so the desktop user actually sees the synced dirs.
MountFlags=shared
ProtectKernelLogs=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
# LIGHTHOUSE-VARMOUNT — create the state dir WITHOUT StateDirectory=. On DO
# droplets /var is a btrfs subvol mounted by initrd, so systemd's generator MASKS
# var.mount; StateDirectory= then auto-adds RequiresMountsFor=/var/lib/mackesd → a
# hard dep on the masked unit → "Failed to start: Unit var.mount is masked" (hit
# live, broke boot-recovery on the lighthouses). An explicit mkdir creates the dir
# (0700, same as StateDirectory) without pulling in any mount dependency, on every
# role. The leading '-' tolerates a pre-existing dir.
ExecStartPre=-/usr/bin/install -d -m 0700 /var/lib/mackesd
TasksMax=4096
MemoryAccounting=true

[Install]
WantedBy=multi-user.target
