# packaging/systemd/magic-mesh-brand.service
[Unit]
Description=MCNF — apply Carbon branding (one-way, first boot, Workstation only)
Documentation=https://github.com/matthewmackes/magic-mesh
# Run after the system is up; the per-user/Cosmic + greeter bits take effect at
# the next login, the system bits (Plymouth, GTK/Qt /etc, LightDM) at next boot.
After=multi-user.target
# BRAND: one-way + run-once. The script writes /var/lib/mde/branded on success.
ConditionPathExists=!/var/lib/mde/branded
# Only on an enrolled, role-pinned node.
ConditionPathExists=/var/lib/mde/role.toml

[Service]
Type=oneshot
RemainAfterExit=yes
# Workstation-role gate: a non-zero ExecCondition SKIPS the unit (not failed),
# so Server/Lighthouse never brand themselves.
ExecCondition=/bin/sh -c 'grep -Eq "role[[:space:]]*=[[:space:]]*\"workstation\"" /var/lib/mde/role.toml'
ExecStart=/usr/libexec/mackesd/magic-mesh-brand apply

[Install]
WantedBy=multi-user.target
