# packaging/systemd/mcnf-farm-autoscale-reconcile.service
# FARM-AUTOSCALE L2 — the elastic-farm reconcile tick, run on the control host
# (fleet-side, NO AI). One oneshot `--once` per fire: read the live per-dom0 build
# demand (worklist @farm jobs + in-flight build slots), decide the per-dom0 shapes
# via install-helpers/farm-autoscale.sh (hysteresis/drain/pods), and — ONLY when
# the apply gate is fully green (FA_APPLY=1 + XO reachable + sane tofu state +
# golden template) — `tofu apply` the shapes. Default is PLAN-ONLY: it never
# applies until the operator opts in AND the prerequisites hold.
#
# Install (operator-gated; NOT enabled by this PR): copy this unit + the timer to
# /etc/systemd/system on the control host, then
#   systemctl enable --now mcnf-farm-autoscale-reconcile.timer
# Leave FA_APPLY unset (plan-only) until XO is up + tofu state is imported; flip it
# with a drop-in (`systemctl edit`) once the live apply prerequisites are met.
[Unit]
Description=MCNF FARM-AUTOSCALE reconciler (demand -> per-dom0 VM shapes; plan-only by default)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=HOME=/root
# Plan-only by default — the apply gate stays closed until the operator opts in.
# Enable the live apply (once XO is up + state imported) via a drop-in:
#   [Service]\nEnvironment=FA_APPLY=1
Environment=FA_APPLY=0
WorkingDirectory=/root/magic-mesh
# MCNF_REPO lets the unit point at the deployed tree; defaults to WorkingDirectory.
ExecStart=/bin/bash -lc '/root/magic-mesh/install-helpers/farm-reconciler.sh --once'
TimeoutStartSec=3600
Nice=10
