# packaging/systemd/mcnf-farm-reconcile.service
# FARM-AUTO-4 — the @farm-build GitOps reconciler, run on the control VM (fleet-
# side, no AI). TEMPLATE; automation/reconciler/reconciler-up.sh installs the live
# copy (EnvironmentFile=/etc/mcnf/reconciler.env, repo slot /opt/mcnf). Never
# mutates infra (no apply gate); dispatches @farm jobs onto the fleet.
[Unit]
Description=MCNF build-farm reconciler (worklist @farm jobs -> fleet builds)
After=network-online.target mcnf-reconciler-bootstrap.service
Wants=network-online.target mcnf-reconciler-bootstrap.service

[Service]
Type=oneshot
Environment=HOME=/root
# Per-mesh config from the rendered env file (no .192/XO literals).
EnvironmentFile=-/etc/mcnf/reconciler.env
WorkingDirectory=/opt/mcnf
ExecStart=/bin/bash -lc 'exec automation/reconciler/farm-reconcile.sh'
TimeoutStartSec=3600
Nice=10
