# packaging/systemd/mcnf-dr-backup.service
# DATACENTER-23 — off-fleet encrypted DR backup of the control plane.
#
# Runs automation/dr/dr-backup.sh, which age-encrypts the OpenTofu state +
# the mesh secret store (+ the Nebula CA when run on a CA holder) into one
# manifest and rclone-pushes it OFF-FLEET to MCNF_DR_REMOTE, so a total
# XCP/DO-fleet loss is recoverable.
#
# NODE PLACEMENT (important): install this on the CONTROL HOST — the node with
# (a) the repo, (b) reach to the LAN datacenter etcd that holds /tofu/state +
# /mcnf/secret, (c) rclone + the `mcnf-spaces` remote configured, and (d) the
# mesh age identity (/root/.mcnf-age-key). A DO lighthouse CANNOT run this
# (it can't reach the LAN datacenter etcd). The Nebula CA itself is separately
# resilient — replicated across every lighthouse; to also push the CA off-fleet,
# run dr-backup.sh once on a CA holder that has rclone + the remote configured.
#
# SAFETY: pushing the CA/secret store off the trust boundary is operator-run by
# design — the safety classifier hard-blocks an automated agent from doing it.
# This unit runs node-local under the operator's systemd, not via any agent.

[Unit]
Description=MCNF DR — off-fleet encrypted control-plane backup (Tofu state + secrets → DO Spaces)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# The off-fleet target + the mesh age identity. Override REPO/REMOTE here if the
# checkout or bucket differ.
Environment=MCNF_DR_REMOTE=mcnf-spaces:mcnf-dr-4533
Environment=MCNF_AGE_KEY=/root/.mcnf-age-key
Environment=MCNF_REPO=/root/magic-mesh
ExecStart=/bin/bash -lc 'cd "${MCNF_REPO}" && exec automation/dr/dr-backup.sh'
# Don't let a transient etcd/rclone hiccup mark the timer failed forever.
SuccessExitStatus=0
