# packaging/systemd/mcnf-dr-backup.service
# DATACENTER-23 / DAR-39/40 — the ON-MESH FIRST LINE of DR (Full tier).
#
# Runs automation/dr/dr-snapshot-onmesh.sh, which produces a v2 dr-<ts>.age
# (dr-backup.sh: tofu state + secrets + age-recipients + Forgejo data via sqlite
# quiesce + a consistent etcd snapshot) and copies it into the Syncthing-replicated
# Mesh-Sync root ($MCNF_MESHFS_DIR/dr/) with INDEX.json + N-deep retention. The
# first recovery line then lands on every surviving peer with NO egress.
#
# The OFF-FLEET push (DO Spaces) is a SEPARATE operator-run step
# (automation/dr/dr-push-offfleet.sh) — NOT fired here. The CA/age master bundle is
# a SEPARATE operator-run passphrase-sealed push (dr-ca-bundle.sh). The safety
# classifier hard-blocks an automated agent from either off-fleet push by design.
#
# NODE PLACEMENT: install on the control VM — it has (a) the repo, (b) reach to the
# founder etcd quorum (via /etc/mackesd/etcd-endpoints — NO .192:2379), (c) the
# Mesh-Sync mount, (d) the mesh age identity. Leader-gated: only the etcd leader's
# run is authoritative (the dr_scheduler worker enforces this; a follower's run is a
# harmless local artifact). Endpoints + paths come from dr-env.sh (DAR-37).

[Unit]
Description=MCNF DR — on-mesh first-line snapshot (v2 manifest → Mesh-Sync)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Per-mesh config (etcd quorum, mesh-storage dir, forgejo data) resolved by
# dr-env.sh; override the repo slot here if the checkout differs.
Environment=MCNF_REPO=/opt/mcnf
EnvironmentFile=-/etc/mcnf/reconciler.env
WorkingDirectory=/opt/mcnf
ExecStart=/bin/bash -lc 'exec automation/dr/dr-snapshot-onmesh.sh'
# Don't let a transient etcd/Mesh-Sync hiccup mark the timer failed forever.
SuccessExitStatus=0
