# packaging/systemd/mde-web-preview-selinux.service
[Unit]
Description=MCNF browser helper SELinux policy setup (non-blocking postinstall)
Documentation=https://github.com/matthewmackes/magic-mesh
After=magic-mesh-selinux-policy.service
Wants=magic-mesh-selinux-policy.service
ConditionPathExists=/usr/libexec/mackesd/setup-selinux-web-preview

[Service]
Type=oneshot
TimeoutStartSec=120
ExecStart=/usr/libexec/mackesd/setup-selinux-web-preview

[Install]
WantedBy=multi-user.target
