# packaging/systemd/mde-musicd.service
[Unit]
Description=Magic Mesh music daemon (Airsonic bridge) — Workstation surfaces only (SVC-7/Q70)
After=network-online.target

[Service]
Type=simple
# SVC-7 — desktop services are Workstation-rank surfaces; a Server or
# Lighthouse box skips this unit cleanly (exit 0 condition failure).
ExecCondition=/usr/bin/mackesd role-gate --min-rank 2
ExecStart=/usr/bin/mde-musicd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
