# packaging/systemd/mesh-status.timer
[Unit]
Description=MCNF — periodic mesh-status snapshot (MESHSHELL SHELL-1)

[Timer]
OnBootSec=20s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
