# packaging/systemd/mesh-shell-setup.service
[Unit]
Description=MCNF — install the starship prompt binary (MESHSHELL SHELL-2)
After=network-online.target
Wants=network-online.target
# Run only until starship is present; retries on later boots if a fetch fails.
ConditionPathExists=!/usr/bin/starship

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/mackesd/mesh-install-starship

[Install]
WantedBy=multi-user.target
