# packaging/systemd/mesh-health.timer
# Drive the MCNF health watchdog ~every minute (first run 90s after boot,
# letting mackesd + nebula settle first).
[Unit]
Description=Periodic MCNF node health watchdog

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
