# packaging/autostart/com.mackes.MagicMeshRoleChooser.desktop
[Desktop Entry]
Type=Application
Name=Construct Role Chooser
Comment=First-run Construct deployment-role chooser; self-exits once a role is pinned, so this is a no-op after the first login.
Exec=mde-role-chooser
Icon=magic-mesh
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
