# packaging/autostart/org.magicmesh.Music.desktop
[Desktop Entry]
Type=Application
Name=MCNF Music Dock
Comment=MUSIC-DOCK-4 — open the MCNF Music bottom dock at login (or as its minimized bottom-center handle). Workstation-rank gated — Servers/Lighthouses are headless, so the role-gate exits before mde-music ever maps a surface.
Exec=sh -c "mackesd role-gate --min-rank 2 && exec mde-music"
Icon=magic-mesh
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
Categories=AudioVideo;Audio;Player;
