# packaging/autostart/org.magicmesh.VoiceAgent.desktop
[Desktop Entry]
Type=Application
Name=Magic Mesh Voice Agent
Comment=Persistent SIP agent (SVC-4); Workstation-rank gated — Servers/Lighthouses run plumbing only (SVC-7/Q70)
Exec=sh -c "mackesd role-gate --min-rank 2 && exec mde-voice-hud --agent"
Icon=magic-mesh
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
